# crontab
# START CRON JOB
* * * * * /usr/local/bin/python3 /app/app-scheduled.py > /proc/1/fd/1 2>/proc/1/fd/2
* * * * * echo “The test cron ran at $(date)” > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB